# tsmicroctl.service.in
[Unit]
Description=tsmicroctl Daemon
After=network.target

[Service]
ExecStart=@bindir@/tsmicroctl --daemon 50
Restart=on-failure
SuccessExitStatus=0
User=root
Group=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
